# service/transitvpc.service
[Unit]
Description=Check S3 bucket:Install Tunnels

[Service]
Type=simple
ExecStart=/opt/transitvpc_ipsec.py

[Install]
WantedBy=multi-user.target
